# deploy/gcp/vacancy-bot.service
[Unit]
Description=Telegram Vacancy Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=/home/REPLACE_USER/vacancy-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/REPLACE_USER/vacancy-bot/.venv/bin/python -m bot.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
